# compilators/debian/nginx/postinst
#!/bin/bash

set -e

if [ "$1" != "configure" ]; then
	exit 0
fi

# Touch and set permisions on default log files on installation
update-rc.d TM defaults > /dev/null
invoke-rc.d TM start || true
